# roles/cd-autoplay/files/playcd.service
[Unit]
Description=Spielt die eingelegte CD

[Service]
Type=oneshot
ExecStart=mplayer -cdrom-device /dev/sr0 cdda:// -cache 5000
