# roles/librespot/files/librespot-watcher.service
[Unit]
Description=Startet Librespot Watcher Script

[Service]
Type=oneshot
ExecStart=/usr/local/bin/librespot-watcher.sh
